# home/systemd/.config/systemd/user/walker.service
[Unit]
Description=Walker launcher daemon
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=/usr/bin/walker --gapplication-service
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
